# config/systemd/user/backup.timer
[Unit]
Description=Daily Backup

[Timer]
OnCalendar=Mon..Fri 12:45:00

[Install]
WantedBy=basic.target
