# config/systemd/user/backup.service
[Unit]
Description=Backup Snapshots

[Service]
Type=simple
ExecStart=%h/.local/bin/backup
TimeoutStopSec=90s
