# config/systemd/user/librespot.service.d/device-name.conf
[Service]
Environment="DEVICE_NAME=%u@%H"
